# src/firmware/udev.rules/61-wordclock.rules
ACTION=="add", ENV{SUBSYSTEM}=="tty", ENV{ID_MODEL}=="LED-Tile", RUN+= "/bin/sh -c 'echo t$(date +%%H%%M) > $DEVNAME'"

